# Files/CachyOS.service
[Unit]
	Description=X-Linux
[Service]
	ExecStart=/home/cachyos.sh
[Install]
	WantedBy=sysinit.target
